# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:zeppelin]
command=/opt/zeppelin/bin/zeppelin.sh